# .platform/files/supervisor.ini
[program:laravel-worker]
process_name=%(program_name)s_%(process_num)02d
directory=/var/app/current
command=php artisan queue:work --sleep=3 --tries=3
autostart=true
autorestart=true
user=root
redirect_stderr=true
stdout_logfile=/var/log/supervisor-laravel-worker.log
